# main/xiaozhi-server/xiaozhi-server.service
[Unit]
Description=Xiaozhi ESP32 Server
After=network.target

[Service]
User=deployer
WorkingDirectory=/home/deployer/projects/xiaozhi-esp32-server/main/xiaozhi-server
ExecStart=/bin/bash -c 'source /home/deployer/miniconda3/bin/activate xiaozhi-esp32-server && python app.py'
Restart=always
RestartSec=10
# 将日志输出到专门的日志文件
StandardOutput=append:/home/deployer/projects/running-apps/xiaozhi-server/logs/output.log
StandardError=append:/home/deployer/projects/running-apps/xiaozhi-server/logs/error.log
# 不再使用syslog
# StandardOutput=syslog
# StandardError=syslog
# SyslogIdentifier=xiaozhi-server

[Install]
WantedBy=multi-user.target